# systemd/user/vulpea.sync.service
[Service]
ExecStart=/usr/bin/env bash -c 'vulpea-sync'
