# systemd/user/vulpea.sync.service
[Service]
ExecStart=/bin/bash -c 'vulpea-sync'
